# supervisord.conf
[unix_http_server]
file = /var/run/supervisor.sock

[supervisord]
nodaemon=true
logfile = /var/log/supervisor/supervisord.log
logfile_maxbytes = 200KB
logfile_backups = 1
pidfile = /var/run/supervisord.pid
childlogdir = /var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

#[program:ovsdb-server]
#command=/usr/sbin/ovsdb-server --remote=punix:/var/run/openvswitch/db.sock --remote=ptcp:6640
#priority=10
#startsec=10
#stderr_events_enabled=true
#stdout_events_enabled=true

#[program:ovs-vswitchd]
#command=/usr/sbin/ovs-vswitchd -v --pidfile
#priority=20
#startsec=10
#stderr_events_enabled=true
#stdout_events_enabled=true

[program:neutron-openvswitch-agent]
command=/usr/bin/neutron-openvswitch-agent --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/plugins/ml2/ml2_conf.ini
priority=30
startsec=10
stderr_events_enabled=true
stdout_events_enabled=true

[program:neutron-server]
command=/usr/bin/neutron-server --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/plugins/ml2/ml2_conf.ini
priority=40
startsec=10
stderr_events_enabled=true
stdout_events_enabled=true
